# launchd/com.renlongyu.rss.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.renlongyu.rss</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/renlongyu/workspace/rss/.venv/bin/python</string>
    <string>-m</string>
    <string>src.main</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/renlongyu/workspace/rss</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/renlongyu/.local/bin</string>
    <key>LANG</key>
    <string>zh_CN.UTF-8</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>10</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/rss-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/rss-daily.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
